# clear_log.crontab
15 3 * * * /apps/sh/act_wdraw_accesslog_cleaner.sh > /dev/null 2>&1
